# warp-cli/warp.conf
[program:warp-svc]
command=/bin/bash -c /usr/bin/warp-svc
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/warp.log
stdout_logfile=/var/log/warp.log
